# dockersh/supervisord-springboot.conf
[program:springboot]
command=/start-springboot.sh
numprocs=1
autostart=true
autorestart=false
user=root
